# deploy/systemd/eldorado-satellite.service
[Unit]
Description=Eldorado Brainrot Scraper (Satellite Node)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/gestion/container/scrapper
EnvironmentFile=/etc/eldorado-scraper/satellite.env
ExecStart=/gestion/container/scrapper/.venv/bin/python /gestion/container/scrapper/scripts/run_dashboard.py
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
